# flowspec/mysql-supervisor.conf
[program:mysqld]
command=/usr/bin/mysqld_safe 
autostart=true
autorestart=true
stderr_logfile=/var/log/mysqld.err.log
stdout_logfile=/var/log/mysqld.out.log
